# vendor/github.com/giantswarm/k8scloudconfig/v_4_10_0/files/conf/10-docker.rules
-w /usr/bin/docker -k docker
-w /var/lib/docker -k docker
-w /etc/docker -k docker
-w /etc/systemd/system/docker.service.d/10-giantswarm-extra-args.conf -k docker
-w /etc/systemd/system/docker.service.d/01-wait-docker.conf -k docker
-w /usr/lib/systemd/system/docker.service -k docker
-w /usr/lib/systemd/system/docker.socket -k docker

