# templates/decrypt-tls-assets.service
[Unit]
Description=Decrypt TLS certificates

[Service]
ExecStart=/opt/bin/decrypt-tls-assets